# docs/io.github.nexustar.usher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent for usher.

  1. Copy to ~/Library/LaunchAgents/ and edit the paths below.
  2. Ensure PATH includes tmux and your CLI (claude / codex).
  3. Load:   launchctl load ~/Library/LaunchAgents/io.github.nexustar.usher.plist
     Unload: launchctl unload ~/Library/LaunchAgents/io.github.nexustar.usher.plist
     Status: launchctl list | grep usher
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.nexustar.usher</string>

  <key>ProgramArguments</key>
  <array>
    <!-- Change to the actual path of your usher binary. -->
    <string>/Users/YOU/go/bin/usher</string>
    <string>serve</string>
    <!-- Add flags here, e.g. <string>--agent-mode</string><string>llm</string> -->
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/usher.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/usher.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- Must include dirs for tmux, claude/codex CLIs, and usher itself. -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/YOU</string>
    <!-- Uncomment for LLM agent mode:
    <key>OPENAI_API_KEY</key>
    <string>sk-xxx</string>
    -->
  </dict>
</dict>
</plist>
